# deploy/rhel/systemd/app.service
[Unit]
Description=PageIndex RAG API Server
After=network.target

[Service]
Type=simple
User=pageindex
Group=pageindex
WorkingDirectory=/opt/pageindex-rag
EnvironmentFile=/opt/pageindex-rag/.env
ExecStart=/opt/pageindex-rag/venv/bin/uvicorn pageindex_rag.api.app:app --host 0.0.0.0 --port 8080
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
